# .debpkg/DEBIAN/prerm
#!/bin/bash
systemctl stop jukebox-player.service

systemctl disable jukebox-player.service

exit 0